# omniorb-4.1.7-cs8-python2/SOURCES/omniNames.service
[Unit]
Description=OmniORB Naming Service
After=syslog.target

[Service]
Type=simple
ExecStart=/usr/bin/omniNames -start -always -logdir /var/log/omniORB/ -errlog /var/log/omniORB/error.log
User=omniORB
Group=omniORB

[Install]
WantedBy=multi-user.target